# ble_orchestrator/systemd/ble-orchestrator.service
[Unit]
Description=BLE Orchestrator Service
After=network.target bluetooth.service

[Service]
ExecStart=/usr/bin/python3 /path/to/ble_orchestrator/main.py
WorkingDirectory=/path/to/ble_orchestrator
Restart=on-failure
StandardOutput=journal
StandardError=journal
TimeoutStopSec=10
RestartSec=5

[Install]
WantedBy=multi-user.target 